# deploy/systemd/ares-cohort-weekly.service
[Unit]
Description=ARES-core recurring-cohort weekly ranking service (port 8000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ares
Group=ares
WorkingDirectory=/opt/ares-core
ExecStart=/opt/ares-core/bin/cohort-service

Environment=COHORT_MODE=weekly
Environment=ARES_HELPER_BINARY=/opt/ares-core/bin/openfhe-helper
Environment=SESSION_PORT=8000
EnvironmentFile=-/etc/ares-core/cohort.env

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ares-core
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
